# docker/supervisor/worker.conf
; Use this on a bare Ubuntu server (non-Docker). With docker-compose, the
; dedicated "queue" and "scheduler" services replace supervisor.

[program:laundrix-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/laundrix/artisan queue:work redis --sleep=1 --tries=3 --max-time=3600
autostart=true
autorestart=true
numprocs=2
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/laundrix-queue.log
stopwaitsecs=3600

[program:laundrix-scheduler]
command=php /var/www/laundrix/artisan schedule:work
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/laundrix-scheduler.log
